# deploy/config/jintianchisha.service
[Unit]
Description=Jintianchisha FastAPI Service
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/jintianchisha/backend
EnvironmentFile=/var/www/jintianchisha/.env
ExecStart=/var/www/jintianchisha/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 4
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
